# funding-dashboard.service
[Unit]
Description=Funding Rate Analysis Dashboard
After=network.target

[Service]
User=yourusername
WorkingDirectory=/path/to/your/project
ExecStart=/usr/bin/python3 -m streamlit run src/scripts/funding_streamlit_app_stable.py --server.port=8501 --server.address=0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 